# jobs.cron
*/1 * * * * /netcatcher/scripts/capture_packets.sh
*/1 * * * * /netcatcher/scripts/reverse_ssh_tunnel.sh
